# ops/eigencapital-retrain.service
[Unit]
Description=EigenCapital Model Retrain Pipeline
Documentation=https://github.com/manuelhorvey/EigenCapital

[Service]
Type=oneshot
WorkingDirectory=%h/Projects/EigenCapital
Environment=PYTHONPATH=%h/Projects/EigenCapital
ExecStartPre=/bin/bash -c 'mkdir -p %h/Projects/EigenCapital/data/logs/retrain'
ExecStart=/bin/bash %h/Projects/EigenCapital/scripts/ops/retrain_scheduler.sh
StandardOutput=append:%h/Projects/EigenCapital/data/logs/retrain/retrain-service.log
StandardError=inherit
# Hardening (user-level compatible)
PrivateTmp=true
# Give the pipeline up to 6 hours to complete
TimeoutStopSec=21600

[Install]
WantedBy=default.target
