# ops/eigencapital-retrain.service
[Unit]
Description=EigenCapital Model Retrain Pipeline
Documentation=https://github.com/manuelhorvey/EigenCapital
Requires=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/Projects/Quorrin
Environment=PYTHONPATH=%h/Projects/Quorrin
ExecStartPre=/bin/bash -c 'mkdir -p %h/Projects/Quorrin/data/logs/retrain'
ExecStart=/bin/bash %h/Projects/Quorrin/scripts/ops/retrain_scheduler.sh
StandardOutput=append:%h/Projects/Quorrin/data/logs/retrain/retrain-service.log
StandardError=inherit
# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
# Give the pipeline up to 6 hours to complete
TimeoutStopSec=21600

[Install]
WantedBy=default.target
